# var/Docker/conf/varnish/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:varnish]
command = bash /start.sh
user = root
autostart = true

